# deploy/systemd/faramesh-serve.service
[Unit]
Description=Faramesh governance daemon (SDK socket + policy)
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
User=faramesh
Group=faramesh
# Adjust paths to your stack (governance.compiled.json from faramesh apply).
ExecStart=/usr/local/bin/faramesh serve --from-compiled /var/lib/faramesh/governance.compiled.json
Restart=always
RestartSec=2
WatchdogSec=30
NotifyAccess=main
RuntimeDirectory=faramesh
RuntimeDirectoryMode=0755
# Hardening (tune for your site)
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
